# scripts/limit.maxfiles.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Raises the per-user file descriptor cap from the default (256 soft /
  unlimited hard on modern macOS — but per-process apps usually inherit
  256) to 65536 soft / 65536 hard. Mutagen + node + Chrome are
  particularly fd-hungry.

  Install ONCE — survives reboot:
    sudo cp scripts/limit.maxfiles.plist /Library/LaunchDaemons/limit.maxfiles.plist
    sudo chown root:wheel /Library/LaunchDaemons/limit.maxfiles.plist
    sudo chmod 644 /Library/LaunchDaemons/limit.maxfiles.plist
    sudo launchctl load -w /Library/LaunchDaemons/limit.maxfiles.plist

  Verify (open a fresh terminal first):
    ulimit -n            # → 65536 soft

  Uninstall:
    sudo launchctl unload -w /Library/LaunchDaemons/limit.maxfiles.plist
    sudo rm /Library/LaunchDaemons/limit.maxfiles.plist
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>limit.maxfiles</string>
  <key>ProgramArguments</key>
  <array>
    <string>launchctl</string>
    <string>limit</string>
    <string>maxfiles</string>
    <string>65536</string>
    <string>65536</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>ServiceIPC</key>
  <false/>
</dict>
</plist>
